# mcproxy.container
[Unit]
Description=MCProxy - MCP Gateway Aggregator
Documentation=https://github.com/mcproxy/mcproxy
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/mcproxy:latest
ContainerName=mcproxy

# Network: Use bridge instead of host for isolation
Network=bridge

# Mount configuration directory (read-only)
Volume=/srv/containers/mcproxy/config:/app/config:ro

# Mount data directory for credentials DB and blocklist cache
Volume=/srv/containers/mcproxy/data:/app/data:rw
Volume=/srv/containers/mcproxy/cache:/app/cache:rw

# Environment variables - pass directly instead of .env file
# This prevents .env from being readable by container
Environment=PYTHONUNBUFFERED=1

# Working directory
WorkingDirectory=/app

# Run with logging to stdout
Exec=python main.py --log --config /app/config/mcp-servers.json

# Container labels
Label=app=mcproxy
Label=service=gateway

# Security settings - HARDENED v4.2
NoNewPrivileges=true
SecurityLabelDisable=false

# Filesystem isolation
ProtectHome=true
ProtectSystem=strict
PrivateTmp=true
ReadOnlyRootfs=true

# Drop all capabilities
CapDrop=ALL

# Resource limits
Memory=512M
MemorySwap=512M

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target